# launchd/com.projectkronos.end-of-day.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.projectkronos.end-of-day</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>__KRONOS_ROOT__/launchd/run-end-of-day.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>900</integer>

  <key>ProcessType</key>
  <string>Background</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>__KRONOS_HOME__</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>StandardOutPath</key>
  <string>__KRONOS_ROOT__/end-of-day.log</string>

  <key>StandardErrorPath</key>
  <string>__KRONOS_ROOT__/end-of-day.log</string>
</dict>
</plist>
